# systemd/user/wssp-daemon.service
[Unit]
Description=Wayland Secret Service Provider (WSSP) Daemon
Documentation=https://github.com/your-repo/wssp
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.secrets
ExecStart=/usr/bin/wssp-daemon
Restart=on-failure
RestartSec=2

# Security hardening
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=read-only
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
RestrictRealtime=yes
UMask=0077

# Allow writing to XDG_RUNTIME_DIR for the IPC socket
ReadWritePaths=%t
# Allow writing to the data directory for the vault
ReadWritePaths=%h/.local/share/wssp

# Ensure GUI environment is available for the prompter
PassEnvironment=WAYLAND_DISPLAY DISPLAY XDG_RUNTIME_DIR XDG_SESSION_ID

# Environment for Headless/IoT deployments (Uncomment to use)
# Environment="WSSP_PASSWORD=your_secure_device_password"

[Install]
WantedBy=default.target
